# private_dot_config/systemd/user/mpd-mpris.service.d/override.conf
[Service]
ExecStart=
ExecStart=mpd-mpris -host="/tmp/mpd_socket" -network=unix
